# plex/supervisord.conf
[supervisord]
nodaemon=true
logfile=/logs/supervisor.log
loglevel=info

[program:plexmediaserver]
command = /opt/plexmediaserver/start_pms
stdout_logfile=/logs/plex.log
